# docs/mcp-deploy/mcp-server.service
# ============================================================
# MCP Server systemd 单元文件（模板）
# install.sh 会替换 __占位符__ 后安装到 /etc/systemd/system/
#
# 设计要点：
#   - EnvironmentFile 指向 .env：改配置只动 .env，不动 service
#   - 用专用用户 mcp 跑，不用 root，最小权限
#   - 安全加固：只读系统、隔离 tmp、禁止提权、限定可写路径
# ============================================================
[Unit]
Description=AI Chat Knowledge MCP Server
Documentation=https://github.com/modelcontextprotocol/python-sdk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__PYTHON_BIN__ __INSTALL_DIR__/mcp_vector_server.py

# 崩溃自动重启
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60

# 标准输出/错误都进 journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-server

# ---- 安全加固 ----
# 禁止通过 setuid 提权
NoNewPrivileges=true

# 系统目录只读（除下方 ReadWritePaths）
ProtectSystem=strict

# 禁止访问 /home
ProtectHome=true

# 隔离临时目录
PrivateTmp=true

# 禁止访问内核日志
ProtectKernelLogs=true

# 禁止访问内核模块
ProtectKernelModules=true

# 禁止修改内核参数
ProtectKernelTunables=true

# 禁止控制组操作
ProtectControlGroups=true

# 限制系统调用到最小集（MCP server 只需网络 + 文件读写）
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# 限制能力（清空所有 capabilities）
CapabilityBoundingSet=
AmbientCapabilities=

# 内存/任务限制（按需调整）
MemoryMax=512M
TasksMax=256

# 只允许写部署目录（日志、缓存等）
ReadWritePaths=__INSTALL_DIR__

[Install]
WantedBy=multi-user.target
